# packages/deb/hazelcast/DEBIAN/hazelcast.service
[Unit]
Description=${HZ_DISTRIBUTION}
After=syslog.target network.target

[Service]
User=hazelcast
Group=hazelcast
Type=simple
ExecStart=hz start
ExecStop=hz-stop
PIDFile=/var/run/hazelcast.pid
StandardOutput=journal

[Install]
WantedBy=multi-user.target
