# packages/deb/hazelcast/DEBIAN/prerm
#!/bin/sh

systemctl stop hazelcast
systemctl disable hazelcast
systemctl daemon-reload

echo "Removing symlinks from /usr/bin"

for FILENAME in /usr/lib/hazelcast/bin/hz*; do
  case "${FILENAME}" in
    *bat)
      ;;
    *)
      rm /usr/bin/"$(basename "${FILENAME}")"
      ;;
  esac
done
